# liams-yoto-sync-tool/rootfs/etc/services.d/liams-yoto-sync-tool/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare config_file_path
declare track_cache_path

## Get the 'config_file_path' key from the user config options.
config_file_path=$(bashio::config 'config_file_path')
track_cache_path=$(bashio::config 'track_cache_path')
auth_store_path=$(bashio::config 'auth_store_path')

## Set dummy environment variables for OAuth2 flows not used by the tool
export SPOTIFY_REDIRECT_URI="DUMMY"
export YOTO_DEV_CLIENT_ID="$(bashio::config 'yoto_dev_client_id')"
export SPOTIFY_CLIENT_ID="$(bashio::config 'spotify_client_id')"

# List media mounts for debugging
bashio::log.info "Listing /media mounts for debugging purposes..."
ls /media

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "${message:="Starting Liam's Yoto Sync Tool..."}"

# Build up arguments
args=("run" "--config" "${config_file_path}")

if bashio::config.has_value 'auth_store_path'; then
    args+=("--auth-store-path" "${auth_store_path}")
fi

if bashio::config.has_value 'track_cache_path'; then
    args+=("--track-cache-path" "${track_cache_path}")
fi

## Run the tool
exec /usr/bin/liams-yoto-sync-tool "${args[@]}"
