# liams-yoto-sync-tool/rootfs/etc/services.d/liams-yoto-sync-tool/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare config_file_path
declare track_cache_path

## Get the 'config_file_path' key from the user config options.
config_file_path=$(bashio::config 'config_file_path')
track_cache_path=$(bashio::config 'track_cache_path')

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "${message:="Starting Liam's Yoto Sync Tool..."}"

## Run the tool
if bashio::config.has_value 'track_cache_path'; then
    exec /usr/bin/liams-yoto-sync-tool run --config "${config_file_path}" --track-cache-path "${track_cache_path}"
else
    exec /usr/bin/liams-yoto-sync-tool run --config "${config_file_path}"
fi
